# packages/thorch-kde-defaults/payload/usr/lib/systemd/user/thorch-audio-init.service
[Unit]
Description=Initialize Thorch speaker audio routing
Wants=pipewire.service wireplumber.service pipewire-pulse.service
After=pipewire.service wireplumber.service pipewire-pulse.service

[Service]
# Audio objects appear asynchronously after PipeWire starts. Do not serialize
# the graphical session behind the routing helper's bounded retry loop.
Type=exec
ExecStart=/usr/bin/thorch-audio-init

[Install]
WantedBy=default.target
