# packages/thorch-bsp/payload/usr/lib/udev/rules.d/90-thorch-backlight.rules
SUBSYSTEM=="backlight", KERNEL=="ae94000.dsi.0", RUN+="/usr/bin/chgrp video /sys/class/backlight/%k/brightness", RUN+="/usr/bin/chmod 0664 /sys/class/backlight/%k/brightness"
SUBSYSTEM=="backlight", KERNEL=="ae96000.dsi.0", RUN+="/usr/bin/chgrp video /sys/class/backlight/%k/brightness", RUN+="/usr/bin/chmod 0664 /sys/class/backlight/%k/brightness"
