# packages/thorch-bsp/payload/usr/lib/systemd/system/thorch-boot-confirm.service
[Unit]
Description=Confirm a successful Thorch boot transaction
After=local-fs.target
ConditionPathExists=/var/lib/thorch/boot-transaction/pending-reboot
ConditionPathIsMountPoint=/boot

[Service]
Type=oneshot
ExecStart=/usr/lib/thorch/boot-transaction confirm

[Install]
WantedBy=multi-user.target
