# packages/thorch-kde-defaults/payload/etc/skel/.config/autostart/thorch-display-setup.desktop
[Desktop Entry]
Type=Application
Name=Thorch Display Setup
Exec=/usr/bin/thorch-display-setup --initial
OnlyShowIn=KDE;
NoDisplay=true
X-KDE-autostart-phase=1
